# systemd/samba_mount/mnt-onedrive.automount
[Unit]
Description=Automount onedrive

[Automount]
Where=/mnt/onedrive

[Install]
WantedBy=multi-user.target
